# docker/supervisor/kafka.conf
[program:kafka]
command=sh /usr/bin/start-kafka.sh
autostart=true
autorestart=false
